# src/deb/init/nntp2p.conf
description "nntp2p usenet proxy"

start on runlevel [2345]
stop on runlevel [!2345]

setuid nntp2p
setgid nntp2p

respawn
exec java -jar /usr/share/nntp2p/nntp2p.jar --settings=/etc/nntp2p.cfg